# root/etc/s6-overlay/s6-rc.d/init-config/run
#!/command/execlineb -P

echo "
The implementation of this image is heavly refered to LinuxServer.io alpine
 base image (https://github.com/linuxserver/docker-baseimage-alpine/tree/master).

To support LSIO projects visit:
https://www.linuxserver.io/donate/
"