# root/etc/s6-overlay/s6-rc.d/init-server-dependency/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "*** Installing required libraries. ***"
python3.11 -m pip install -r /app/requirements.txt